# SendToAio/SendToAio.crontab
#
#
*/15 * * * * /home/pi/SendToAio/run.sh >> /var/log/SendToAio 2>&1
